# backend/deploy/havalielaletlerisatis.service
# /etc/systemd/system/havalielaletlerisatis.service
[Unit]
Description=Havalı El Aletleri Satış FastAPI
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/havalielaletlerisatis/backend/ecommerce_app
EnvironmentFile=/var/www/havalielaletlerisatis/backend/ecommerce_app/.env
ExecStart=/var/www/havalielaletlerisatis/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
